# rootfs/etc/services.d/sshd/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh:
/usr/sbin/sshd -D -f /etc/ssh/sshd_config
