# rootfs/etc/services.d/karaf/run
#!/usr/bin/execlineb -P
# -*- mode: sh -*-
# vi: set ft=sh:
# We should run karaf as it's own user, but it just doesn't work and I've spent
# way to much time on this (this is security issue).
# s6-setuidgid karaf
with-contenv
multisubstitute
{
    import PATH
    import JAVA_HOME
    import JAVA_OPTS
}
s6-env -i
PATH=$PATH
JAVA_HOME=$JAVA_HOME
JAVA_OPTS=$JAVA_OPTS
/opt/karaf/bin/karaf server
